# Firmware/src/qml-plugin/tools/cellulorobothubd/scripts/cellulorobothubd.service
[Unit]
Description=Cellulo Robot Hub Daemon
After=network.target

[Service]
Type=forking
ExecStart=/usr/local/bin/cellulorobothubd

[Install]
WantedBy=multi-user.target
